# deploy/supervisor.conf
# ════════════════════════════════════════════════════════════════
# Supervisor конфигурации для your-domain.kz
# Файл: /etc/supervisor/conf.d/chairs-almaty.conf
# После изменений: supervisorctl reread && supervisorctl update
# ════════════════════════════════════════════════════════════════

# ── Queue Worker: обычные задачи ──────────────────────────────
[program:chairs-almaty-worker]
command=php /var/www/your-domain.kz/artisan queue:work database
    --queue=imports-high,imports,imports-low,default
    --sleep=3
    --tries=3
    --timeout=120
    --max-time=3600
    --max-jobs=1000
process_name=%(program_name)s_%(process_num)02d
numprocs=2
directory=/var/www/your-domain.kz
user=www-data
autostart=true
autorestart=true
startretries=3
startsecs=3
stdout_logfile=/var/log/supervisor/chairs-almaty-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/chairs-almaty-worker-error.log

# ── Scheduler: cron-задачи Laravel ───────────────────────────
[program:chairs-almaty-scheduler]
command=php /var/www/your-domain.kz/artisan schedule:work
process_name=%(program_name)s
numprocs=1
directory=/var/www/your-domain.kz
user=www-data
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/chairs-almaty-scheduler.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/supervisor/chairs-almaty-scheduler-error.log
